# crontab
*/10 * * * * root cd /root/lewdcrops && ./main.py >> /var/log/cron.log 2>&1
